# launchd/dev.influo.mlx-lazyserve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.influo.mlx-lazyserve</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/uv</string>
        <string>run</string>
        <string>--no-sync</string>
        <string>--directory</string>
        <string>~/mlx-lazyserve</string>
        <string>mlx-lazyserve</string>
    </array>

    <!-- Bind on all interfaces so the Tailscale IP is reachable. Prefer setting
         the Tailscale 100.x address here, or set MLX_LAZYSERVE_API_KEYS. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>MLX_LAZYSERVE_HOST</key>
        <string>0.0.0.0</string>
        <key>MLX_LAZYSERVE_PORT</key>
        <string>41434</string>
        <key>MLX_LAZYSERVE_IDLE_TIMEOUT</key>
        <string>600</string>
        <!-- Raise the Metal GPU wired-memory limit on start (reset to 0 on stop) so the
             ~18.7 GB Qwen3.6-A3B fits past the measured ~17.8 GB default cap.
             Needs the NOPASSWD sudoers rule (launchd/mlx-lazyserve.sudoers). -->
        <key>MLX_LAZYSERVE_WIRED_LIMIT_MB</key>
        <string>22000</string>
        <!-- These repos are xet-backed; xet transfer is slow from CN (~0.5 MB/s).
             Disabling xet uses the classic HF CDN path (~14 MB/s) with auto-resume. -->
        <key>HF_HUB_DISABLE_XET</key>
        <string>1</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <!-- Interactive avoids macOS background throttling so inference runs full-speed. -->
    <key>ProcessType</key>
    <string>Interactive</string>

    <key>StandardOutPath</key>
    <string>~/mlx-lazyserve/logs/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>~/mlx-lazyserve/logs/stderr.log</string>
</dict>
</plist>
